# files/rsc-format-disk@.service
[Unit]
Description=Formats Research Cloud Disk on %i

[Service]
Type=forking
GuessMainPID=no
ExecStart=/opt/rsc-utilities/format_rsc_disk.sh %I